# photo-booth.service
[Unit]
Description=Photo Booth Camera Service
Documentation=https://github.com/yourusername/photo-booth
After=network.target cups.service
Wants=network-online.target

[Service]
Type=simple
User=prospero
Group=prospero
WorkingDirectory=/home/prospero

# Environment variables
Environment="VIDEO_DEVICE=/dev/video0"
Environment="VIDEO_WIDTH=1920"
Environment="VIDEO_HEIGHT=1080"
Environment="VIDEO_FORMAT=MJPG"
Environment="STORAGE_PATH=/usr/local/share/photo_booth"
Environment="PRINTER_NAME=XP8700series-TurboPrint"
Environment="PRINTER_FALLBACK=EPSON_XP_8700_Series_USB,XP-8700"
Environment="USE_MOCK_PRINTER=false"
Environment="TEMPLATE_HEADER=ESSEN SPIEL '25"
Environment="TEMPLATE_NAME=NO GOODUN'"
Environment="TEMPLATE_HEADLINE=RUSTLIN' & BEWITCHIN' CATTLE
Environment="TEMPLATE_STORY=Once there was a little place with a lot to do and not a lot to think. Then there came a place with naught to feel, but terror in the day. Then there came so much to feel and oh too much to think."
Environment="TEMPLATE_BACKGROUND=background.png"
Environment="RUST_LOG=info"
Environment="HOST=0.0.0.0"
Environment="PORT=8080"

# Load additional environment from file if it exists
EnvironmentFile=-/home/prospero/.env

# Execute the binary
ExecStart=/home/prospero/cam_test

# Restart policy
Restart=always
RestartSec=10

# Resource limits
LimitNOFILE=4096
MemoryLimit=512M

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/usr/local/share/photo_booth

# Supplementary groups for device access
SupplementaryGroups=video dialout lp

# Standard output/error logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=photo-booth

[Install]
WantedBy=multi-user.target
